# contrib/systemd/agent-diva.service
[Unit]
Description=Agent Diva Gateway Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=agentdiva
Group=agentdiva
WorkingDirectory=/opt/agent-diva
ExecStart=/usr/bin/agent-diva gateway
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agent-diva /var/log/agent-diva
PrivateTmp=true

# 环境变量
Environment="RUST_LOG=info"
Environment="AGENT_DIVA_CONFIG=/opt/agent-diva/config.json"

[Install]
WantedBy=multi-user.target